# rootfs/etc/init.d/coyote-config
#!/sbin/openrc-run
#
# Coyote Linux configuration application service
#
# This service applies the Coyote configuration including network setup.
# It provides the 'networking' service so other services can depend on it.
#

description="Apply Coyote Linux configuration"

depend() {
    need coyote-init
    provide net
    before coyote-webadmin
}

start() {
    ebegin "Applying Coyote Linux configuration"

    # Apply system configuration (includes network setup, mounts config partition)
    /opt/coyote/bin/apply-config

    # Set up persistent directories now that config partition is mounted
    # /root is a symlink to /mnt/config/root for persistent home
    setup_persistent_dirs

    # Set up web admin write permissions
    setup_webadmin_permissions

    eend $?
}

# Create persistent directories on config partition
setup_persistent_dirs() {
    # Check if config partition is mounted
    if ! mountpoint -q /mnt/config 2>/dev/null; then
        ewarn "Config partition not mounted, skipping persistent directory setup"
        return 0
    fi

    local need_write=false

    # Check if directories need to be created
    [ ! -d /mnt/config/root ] && need_write=true
    [ ! -d /mnt/config/ssl ] && need_write=true

    if [ "$need_write" = "true" ]; then
        einfo "Creating persistent directories..."
        mount -o remount,rw /mnt/config 2>/dev/null || {
            ewarn "Cannot remount config partition read-write"
            return 0
        }

        # Create root home directory
        if [ ! -d /mnt/config/root ]; then
            mkdir -p /mnt/config/root
            chmod 700 /mnt/config/root
        fi

        # Create .ssh directory for root
        if [ ! -d /mnt/config/root/.ssh ]; then
            mkdir -p /mnt/config/root/.ssh
            chmod 700 /mnt/config/root/.ssh
        fi

        # Create SSL certificate directory
        if [ ! -d /mnt/config/ssl ]; then
            mkdir -p /mnt/config/ssl
            chmod 700 /mnt/config/ssl
        fi

        mount -o remount,ro /mnt/config 2>/dev/null || true
    fi
}

# Set up permissions for web admin to write configuration
setup_webadmin_permissions() {
    # Create running-config directory with lighttpd ownership
    # This holds the configuration currently applied to the system
    if [ ! -d /tmp/running-config ]; then
        mkdir -p /tmp/running-config
    fi
    chown lighttpd:lighttpd /tmp/running-config
    chmod 755 /tmp/running-config

    # Create working-config directory with lighttpd ownership
    # This holds uncommitted changes made in the web admin
    if [ ! -d /tmp/working-config ]; then
        mkdir -p /tmp/working-config
    fi
    chown lighttpd:lighttpd /tmp/working-config
    chmod 755 /tmp/working-config

    # If running-config file exists (created by apply-config), make it writable
    if [ -f /tmp/running-config/system.json ]; then
        chown lighttpd:lighttpd /tmp/running-config/system.json
        chmod 644 /tmp/running-config/system.json
    fi

    # Set permissions on persistent config for web admin writes
    # Config partition must be mounted
    if mountpoint -q /mnt/config 2>/dev/null; then
        # Temporarily remount rw to fix permissions
        mount -o remount,rw /mnt/config 2>/dev/null || return 0

        # Make system.json writable by lighttpd if it exists
        if [ -f /mnt/config/system.json ]; then
            chown lighttpd:lighttpd /mnt/config/system.json
            chmod 644 /mnt/config/system.json
        fi

        # Create and set permissions on backups directory
        if [ ! -d /mnt/config/backups ]; then
            mkdir -p /mnt/config/backups
        fi
        chown lighttpd:lighttpd /mnt/config/backups
        chmod 755 /mnt/config/backups

        mount -o remount,ro /mnt/config 2>/dev/null || true
    fi
}

stop() {
    ebegin "Stopping Coyote configuration"
    eend 0
}

reload() {
    ebegin "Reloading Coyote configuration"
    /opt/coyote/bin/apply-config
    eend $?
}
