# rootfs/etc/init.d/coyote-config
#!/sbin/openrc-run
#
# Coyote Linux configuration application service
#
# This service applies the Coyote configuration including network setup.
# It provides the 'networking' service so other services can depend on it.
#

description="Apply Coyote Linux configuration"

depend() {
    need coyote-init
    provide net
    before coyote-webadmin
}

start() {
    ebegin "Applying Coyote Linux configuration"

    # Apply system configuration (includes network setup)
    /opt/coyote/bin/apply-config

    eend $?
}

stop() {
    ebegin "Stopping Coyote configuration"
    eend 0
}

reload() {
    ebegin "Reloading Coyote configuration"
    /opt/coyote/bin/apply-config
    eend $?
}
