# rootfs/etc/init.d/coyote-webadmin
#!/sbin/openrc-run
#
# Coyote Linux Web Administration Interface
#

description="Coyote Web Administration Interface"
name="coyote-webadmin"

depend() {
    need coyote-config
    after networking
}

start_pre() {
    # lighttpd user/group created during build - verify it exists
    if ! id -u lighttpd >/dev/null 2>&1; then
        eerror "lighttpd user does not exist"
        return 1
    fi

    # Ensure log directory exists
    checkpath --directory --owner lighttpd:lighttpd --mode 0755 /var/log/lighttpd

    # Ensure socket directory exists
    checkpath --directory --owner lighttpd:lighttpd --mode 0755 /var/run/lighttpd

    # SSL certificates stored on config partition for persistence
    # Create symlink from /etc/lighttpd/ssl to /mnt/config/ssl
    mkdir -p /etc/lighttpd
    if [ ! -L /etc/lighttpd/ssl ]; then
        rm -rf /etc/lighttpd/ssl
        ln -s /mnt/config/ssl /etc/lighttpd/ssl
    fi

    # Generate SSL certificate if it doesn't exist (stored persistently)
    if [ ! -f /mnt/config/ssl/server.pem ]; then
        einfo "Generating self-signed SSL certificate..."

        # Remount config partition read-write
        mount -o remount,rw /mnt/config

        # Ensure SSL directory exists
        mkdir -p /mnt/config/ssl
        chmod 700 /mnt/config/ssl

        # Generate certificate
        openssl req -new -x509 -days 3650 -nodes \
            -out /mnt/config/ssl/server.pem \
            -keyout /mnt/config/ssl/server.pem \
            -subj "/CN=$(hostname)/O=Coyote Linux/C=US" 2>/dev/null
        chmod 600 /mnt/config/ssl/server.pem

        # Remount config partition read-only
        mount -o remount,ro /mnt/config
    fi
}

start() {
    ebegin "Starting Coyote Web Admin"
    start-stop-daemon --start --exec /usr/sbin/lighttpd -- -f /etc/lighttpd/lighttpd.conf
    eend $?
}

stop() {
    ebegin "Stopping Coyote Web Admin"
    start-stop-daemon --stop --exec /usr/sbin/lighttpd
    eend $?
}

reload() {
    ebegin "Reloading Coyote Web Admin"
    start-stop-daemon --stop --signal HUP --exec /usr/sbin/lighttpd
    eend $?
}
