# rootfs/etc/init.d/coyote-config
#!/sbin/openrc-run
#
# Coyote Linux configuration application service
#

description="Apply Coyote Linux configuration"

depend() {
    need coyote-init
    need network
}

start() {
    ebegin "Applying Coyote Linux configuration"

    # Apply system configuration
    /opt/coyote/bin/apply-config

    eend $?
}

stop() {
    ebegin "Stopping Coyote configuration"
    eend 0
}

reload() {
    ebegin "Reloading Coyote configuration"
    /opt/coyote/bin/apply-config
    eend $?
}
