# rootfs/etc/init.d/coyote-webadmin
#!/sbin/openrc-run
#
# Coyote Linux Web Administration Interface
#

description="Coyote Web Administration Interface"
name="coyote-webadmin"

depend() {
    need net
    need coyote-config
    after firewall
}

start_pre() {
    # Ensure lighttpd user exists (created by Alpine lighttpd package)
    if ! id -u lighttpd >/dev/null 2>&1; then
        addgroup -S lighttpd 2>/dev/null
        adduser -S -D -H -h /var/www -s /sbin/nologin -G lighttpd lighttpd 2>/dev/null
    fi

    # Ensure log directory exists
    checkpath --directory --owner lighttpd:lighttpd --mode 0755 /var/log/lighttpd

    # Ensure socket directory exists
    checkpath --directory --owner lighttpd:lighttpd --mode 0755 /var/run/lighttpd

    # Ensure SSL directory exists
    mkdir -p /etc/lighttpd/ssl

    # Ensure SSL certificate exists
    if [ ! -f /etc/lighttpd/ssl/server.pem ]; then
        einfo "Generating self-signed SSL certificate..."
        openssl req -new -x509 -days 3650 -nodes \
            -out /etc/lighttpd/ssl/server.pem \
            -keyout /etc/lighttpd/ssl/server.pem \
            -subj "/CN=$(hostname)/O=Coyote Linux/C=US" 2>/dev/null
        chown lighttpd:lighttpd /etc/lighttpd/ssl/server.pem
        chmod 600 /etc/lighttpd/ssl/server.pem
    fi
}

start() {
    ebegin "Starting Coyote Web Admin"
    start-stop-daemon --start --exec /usr/sbin/lighttpd -- -f /etc/lighttpd/lighttpd.conf
    eend $?
}

stop() {
    ebegin "Stopping Coyote Web Admin"
    start-stop-daemon --stop --exec /usr/sbin/lighttpd
    eend $?
}

reload() {
    ebegin "Reloading Coyote Web Admin"
    start-stop-daemon --stop --signal HUP --exec /usr/sbin/lighttpd
    eend $?
}
