# rootfs/etc/init.d/root
#!/sbin/openrc-run
#
# Coyote Linux root filesystem handler
#
# This replaces the default Alpine root script which attempts to remount
# the root filesystem as read-write. Coyote uses a read-only squashfs
# root filesystem by design, so no remount is needed or possible.
#

description="Coyote Linux root filesystem (read-only squashfs)"

depend() {
    # No dependencies - this is one of the first services to run
    # We explicitly do NOT need localmount since our root is read-only squashfs
    provide root
    before localmount
}

start() {
    # Nothing to do - root is already mounted read-only by initramfs
    ebegin "Root filesystem (read-only squashfs)"
    eend 0
}

stop() {
    return 0
}
